# overlay/etc/init.d/20_bluetooth
#!/sbin/openrc-run

command="/etc/init.d/bluetooth start || true"