# overlay/etc/init.d/20_wifi
#!/sbin/openrc-run

description="ESP32 SPI hosted driver"

depend() {
    need localmount
}

start() {
    ebegin "Loading esp32_spi module"
    insmod /var/esp32_spi.ko resetpin=54 clockspeed=30 || eend 1

    eend 0
}

stop() {
    ebegin "Unloading esp32_spi module"
    rmmod esp32_spi 2>/dev/null

    eend 0
}
