# overlay/etc/init.d/30_wifi
#!/sbin/openrc-run

command="/usr/bin/wifi_service start || true"